# deploy/pytombo.service
[Unit]
Description=TTombo (pytombo) Streamlit app
After=network.target

[Service]
User=pytombo
WorkingDirectory=/opt/pytombo
ExecStart=/opt/pytombo/.venv/bin/streamlit run app.py --server.port 8501 --server.address 127.0.0.1 --server.headless true --browser.gatherUsageStats false
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
